# pc-huawei/x-watchdog.service
[Unit]
Description=X Monitor Watchdog v2.1 — 服务健康监控 + 自愈
After=network.target x-monitor.service x-web.service
Wants=x-monitor.service x-web.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/x-monitor
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="LANG=en_US.UTF-8"
ExecStart=/opt/x-monitor/venv/bin/python /opt/x-monitor/watchdog.py --daemon --interval 300
Restart=always
RestartSec=30
StandardOutput=append:/opt/x-monitor/watchdog.log
StandardError=append:/opt/x-monitor/watchdog.log

[Install]
WantedBy=multi-user.target
